# ros_files/devel/etc/catkin/profile.d/25.turtlebot.sh
/home/petar/project_ws/devel/.private/turtlebot_bringup/etc/catkin/profile.d/25.turtlebot.sh